# deploy/massive-backfill.timer
[Unit]
Description=Run Options Flywheel Massive backfill chunks

[Timer]
OnBootSec=5min
OnUnitInactiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
