# deploy/sata-options.service
[Unit]
Description=Options Flywheel
After=network.target

[Service]
Type=simple
User=main
Group=main
WorkingDirectory=/home/main/options
Environment=APP_BASE_PATH=/options
Environment=MARKET_DATA_PROVIDER=yahoo
Environment=DATABASE_URL=sqlite:///./sata_options_optimizer.db
ExecStart=/home/main/options/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
